# bash/profile.d/more_execution_folders.sh
export PATH=$HOME/mybin:$PATH
export PATH=$HOME/Code/scripts:$PATH

